# root/etc/NetworkManager/dispatcher.d/10-tzupdate.sh
#!/bin/sh
status=$2
case $status in
    up)
        /usr/bin/tzupdate
esac
