# root/etc/systemd/system/ipmi/ipmi-temp.timer
[Unit]
Description=Watches server temperature

[Timer]
OnUnitInactiveSec=5s
OnBootSec=5s

[Install]
WantedBy=timers.target
