# root/etc/systemd/system/ipmi-static.service
[Unit]
Description=Set fan speed

[Service]
Type=simple
ExecStart=/home/balrog/Scripts/ipmi-static.sh

[Install]
WantedBy=multi-user.target
