# root/etc/systemd/system/nightly-reboot.service
[Unit]
Description=Nightly reboot

[Service]
Type=simple
ExecStart=/usr/bin/systemctl reboot
