# root/etc/systemd/system/pacman-dl.timer
[Unit]
Description=Download packages hourly

[Timer]
OnCalendar=hourly
