# home/.config/systemd/user/libra.service
[Unit]
Description=Listenbrainz scrobbler for MPD
After=mpd.service

[Service]
ExecStart=%h/.bin/go/bin/libra

[Install]
WantedBy=default.target
