# home/.config/systemd/user/mpDris2.service
[Unit]
Description=MPRIS2 support for MPD
After=mpd.service

[Service]
Type=simple
ExecStart=/usr/bin/mpDris2

[Install]
WantedBy=default.target
