# home/.config/systemd/user/trash-empty.timer
[Unit]
Description=Runs trash-empty.service every hour

[Timer]
OnUnitActiveSec=1h
Unit=trash-empty.service

[Install]
WantedBy=timers.target
