# home/.config/systemd/user/trash-empty.service
[Unit]
Description=Removes files older than a week from trash can

[Service]
Type=oneshot
ExecStart=/usr/bin/trash-empty 7
